# senses-ecom-bi/deploy/senses-ingest-nuvemshop.service
[Unit]
Description=Senses BI — ingest Nuvemshop (financeiro) via container
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/root/senses/senses-ecom-bi
# Roda o ingest dentro do docker (rede supabase_default → alcança supabase-db).
# Janela: 3 dias atrás → HOJE (inclui o dia corrente p/ atualizar quase ao vivo).
# Upsert idempotente: re-rodar de 30 em 30 min só sobrescreve os mesmos dias.
ExecStart=/bin/sh -c '/usr/bin/docker compose run --rm ingest ingest_nuvemshop.py --from "$(date -d "-3 days" +%%F)" --to "$(date +%%F)"'
StandardOutput=journal
StandardError=journal
