# senses-ecom-bi/deploy/senses-ingest-nuvemshop.timer
[Unit]
Description=Senses BI — ingest Nuvemshop diário (~06:00, após fechamento do dia)

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
