# senses-ecom-bi/deploy/senses-ingest-nuvemshop.timer
[Unit]
Description=Senses BI — ingest Nuvemshop a cada 30 min (dia de hoje quase ao vivo)

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
