# senses-ecom-bi/deploy/senses-ingest-ga4.timer
[Unit]
Description=Senses BI — ingest GA4 diário (~06:05, após o financeiro)

[Timer]
OnCalendar=*-*-* 06:05:00
Persistent=true

[Install]
WantedBy=timers.target
