# field-dashboard/sync/systemd/field-sync.timer
[Unit]
Description=Roda o sync Field->Supabase a cada 30 minutos

[Timer]
OnCalendar=*:0/30
Persistent=true
RandomizedDelaySec=60
Unit=field-sync.service

[Install]
WantedBy=timers.target
